# security/tor-devel/files/tor.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: tor
# REQUIRE: DAEMON FILESYSTEMS
# BEFORE:  LOGIN
#
# tor_enable	(bool):	Set it to "YES" to enable Tor.	Default: NO
# tor_conf	(str):	Tor's config file.		Default: %%PREFIX%%/etc/tor/torrc
# tor_pidfile	(str):	Tor's PID file.			Default: /var/run/tor/tor.pid
#
#   This options are command line options used by Tor.
#     That said, setting up 'DataDirectory' or 'Log' using these 
#     options will result in configuration overwriting.
#
# tor_datadir	(str):	Tor's DataDirecoty.		Default: /var/db/tor
# tor_flags	(str):	Extra args used to run Tor.	Default: --Log ${tor_loglevel}\\ file\\ ${tor_logdir}/${tor_loglevel}.log --RunAsDaemon 1
# tor_logdir	(str):	Log directory.			Default: /var/log/tor
# tor_loglevel	(str):	Log severity level.		Default: notice
# tor_setuid	(bool):	Run Tor with setuid bit on.	Default: NO
# tor_user	(str):	User to run the Tor daemon.	Default: %%USER%%
#

. /etc/rc.subr

name="tor"
rcvar="${name}_enable"

load_rc_config ${name}

: ${tor_enable="NO"}
: ${tor_conf="%%PREFIX%%/etc/tor/torrc"}
: ${tor_pidfile="/var/run/tor/tor.pid"}

: ${tor_datadir="/var/db/tor"}
: ${tor_flags="--Log ${tor_loglevel}\\ file\\ ${tor_logdir}/${tor_loglevel}.log --RunAsDaemon 1"}
: ${tor_logdir="/var/log/tor"}
: ${tor_loglevel="notice"}
: ${tor_setuid="NO"}
: ${tor_user="%%USER%%"}

checkyesno tor_setuid && tor_user="root"

required_files="${tor_conf}"
required_dirs="${tor_datadir}"

pidfile="${tor_pidfile}"

command="%%PREFIX%%/bin/tor"
command_args="-f ${tor_conf} -User %%USER%% --DataDirectory ${tor_datadir} --PidFile ${tor_pidfile} ${tor_flags}"

extra_commands="reload"

run_rc_command $1
