# security/tor-devel/files/tor.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: tor
# REQUIRE: DAEMON FILESYSTEMS
# BEFORE: LOGIN
#
# tor_enable (bool):	Set it to "YES" to enable Tor.	Default: NO
# tor_conf (str):	Tor's config file.		Default: %%PREFIX%%/etc/tor/torrc
#
# This options are command line options used by Tor. That said,
# 	setting up 'DataDirectory' or 'Log' using these 
# 	options will result in configuration overwriting.
#
# tor_user (str):	User to run the Tor daemon.	Default: %%USER%%
# tor_datadir (str):	Tor's DataDirecoty.		Default: /var/db/tor
# tor_logdir (str):	Log directory.			Default: /var/log/tor
# tor_loglevel (str):	Log severity level.		Default: notice
#

. /etc/rc.subr

name="tor"
rcvar=tor_enable

load_rc_config ${name}

: ${tor_enable="NO"}
: ${tor_conf="%%PREFIX%%/etc/tor/torrc"}
: ${tor_user="%%USER%%"}
: ${tor_pidfile="/var/run/tor/tor.pid"}
: ${tor_datadir="/var/db/tor"}
: ${tor_logdir="/var/log/tor"}
: ${tor_loglevel="notice"}

required_files=${tor_conf}
required_dirs=${tor_datadir}
pidfile=${tor_pidfile}
command="%%PREFIX%%/bin/${name}"
command_args="-f ${tor_conf} --RunAsDaemon 1 --PidFile ${tor_pidfile} --DataDirectory ${tor_datadir} --Log ${tor_loglevel}\\ file\\ ${tor_logdir}/${tor_loglevel}.log"
extra_commands="reload"

run_rc_command "$1"

