# dns/dnscrypt-proxy2/files/dnscrypt-proxy.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: dnscrypt_proxy
# REQUIRE: SERVERS cleanvar ldconfig
# BEFORE: dnsmasq local_unbound unbound named
# KEYWORD: shutdown
#
# dnscrypt_proxy_enable	(bool):	Enable service on boot.	Default: NO
# dnscrypt_proxy_config	(str):	Config file.		Default: %%PREFIX%%/etc/dnscrypt-proxy.toml
# dnscrypt_proxy_flags	(str):	Additional flags.	
#	Default: -config ${dnscrypt_proxy_config} -logfile /var/log/${name}.log
#

. /etc/rc.subr

name="dnscrypt_proxy"
rcvar="${name}_enable"

load_rc_config "${name}"

: ${dnscrypt_proxy_enable:="NO"}
: ${dnscrypt_proxy_config:="%%PREFIX%%/etc/dnscrypt-proxy.toml"}
: ${dnscrypt_proxy_command:="%%PREFIX%%/sbin/dnscrypt-proxy"}
: ${dnscrypt_proxy_flags:="-config ${dnscrypt_proxy_config} -logfile /var/log/${name}.log"}
: ${dnscrypt_proxy_user:="_dnscrypt-proxy"}
: ${dnscrypt_proxy_pidfile:="/var/run/dnscrypt-proxy.pid"}

pidfile="${dnscrypt_proxy_pidfile}"
required_files="${dnscrypt_proxy_config}"

command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -f ${dnscrypt_proxy_command}"

run_rc_command $1
