# dns/dnscrypt-proxy/files/dnscrypt-proxy.in
#!/bin/sh
#
# $FreeBSD: head/dns/dnscrypt-proxy/files/dnscrypt-proxy.in 397925 2015-09-26 02:39:09Z junovitch $
#
# PROVIDE: dnscrypt_proxy
# REQUIRE: SERVERS cleanvar ldconfig
# BEFORE: dnsmasq named local_unbound unbound pdnsd
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable dnscrypt-proxy:
#
# dnscrypt_proxy_enable (bool):	Set to NO by default.
#				Set to YES to enable dnscrypt-proxy.
# dnscrypt_proxy_uid (str):	Set to "_dnscrypt-proxy" by default.
#                              	User to switch to after starting.
# dnscrypt_proxy_resolver (str):Set to "dnscrypt.eu-nl" by default.
#                              	Choose a different upstream resolver.
# dnscrypt_proxy_pidfile (str):	default: "/var/run/dnscrypt-proxy.pid"
#                              	Location of pid file.
# dnscrypt_proxy_logfile (str):	default: "/var/log/dnscrypt-proxy.log"
#                              	Location of log file.
# dnscrypt_proxy_flags (str):	Custom flags to use with dnscrypt-proxy

. /etc/rc.subr

name=dnscrypt_proxy
rcvar=dnscrypt_proxy_enable

load_rc_config ${name}

: ${dnscrypt_proxy_enable:=NO}
: ${dnscrypt_proxy_uid=_dnscrypt-proxy}			# User to run daemon as
: ${dnscrypt_proxy_resolver=dnscrypt.eu-nl}		# Resolver to be used
: ${dnscrypt_proxy_pidfile=/var/run/dnscrypt-proxy.pid} # Path to pid file
: ${dnscrypt_proxy_logfile=/var/log/dnscrypt-proxy.log} # Path to log file
: ${dnscrypt_proxy_flags=""}				# Custom flags to use

command=%%PREFIX%%/sbin/dnscrypt-proxy
command_args="-d -p ${dnscrypt_proxy_pidfile} -l ${dnscrypt_proxy_logfile} -u ${dnscrypt_proxy_uid} -R ${dnscrypt_proxy_resolver} ${dnscrypt_proxy_flags}"
procname=%%PREFIX%%/sbin/dnscrypt-proxy
pidfile=${dnscrypt_proxy_pidfile}

run_rc_command "$1"
